# aggregator/systemd/ai-spend-http.service
[Unit]
Description=AI Spend JSON HTTP Server
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 -m http.server 9876 --bind 0.0.0.0
WorkingDirectory=/opt/hermes/public
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
